# etc/profile.d/dns.sh
PATH=${PATH}:/opt/dns/bin:/opt/dns/sbin
alias ddns="cd /opt/dns"
